# ansible/roles/system/files/fancontrol.service
[Unit]
Description=Fancontrol service

[Service]
ExecStart=/usr/bin/fancontrol.py

[Install]
WantedBy=multi-user.target
