# hack/ceph-backup.timer
[Unit]
Description=Backup Ceph Daily

[Timer]
OnCalendar=daily
AccuracySec=1m

[Install]
WantedBy=multi-user.target
